# scripts/yinsh-dashboard.service
[Unit]
Description=YinshML Metrics Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=yinsh
Group=yinsh
WorkingDirectory=/opt/yinshml
Environment=PATH=/opt/yinshml/venv/bin
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=false
Environment=SECRET_KEY=your-production-secret-key-here
Environment=JWT_SECRET_KEY=your-jwt-secret-key-here
Environment=DEFAULT_ADMIN_PASSWORD=your-admin-password-here
ExecStart=/opt/yinshml/venv/bin/python dashboard/app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yinsh-dashboard

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/yinshml
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target 